# crontab
0 5 * * * ./blocksite.sh enable
0 18 * * * ./blocksite.sh disable